# service/tiramisu/run
#!/bin/sh
pgrep tiramisu || tiramisu -o '#summary #body' | \
	while read -r line; do
		# Make sure notification stays for 10 sec
		for i in $(seq 10); do
			herbe "$line"
			sleep 1
		done
	done
